# embedded-evnt-controller/deployment/yocto/meta-vehicle/recipe-apps/vehicle-controller/file/vehicle-controller.service
[Unit]
Description=Vehicle Controller Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/vehicle-controller
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target